# packaging/freebsd/rc.d/openwatch_worker
#!/bin/sh
#
# PROVIDE: openwatch_worker
# REQUIRE: LOGIN postgresql openwatch_api
# KEYWORD: shutdown
#
# OpenWatch background worker service (PostgreSQL-backed job queue)
#
# Add the following lines to /etc/rc.conf to enable:
#   openwatch_worker_enable="YES"
#
# Optional rc.conf settings:
#   openwatch_worker_user="openwatch"  # Run as user (default: openwatch)
#   openwatch_worker_logfile="/var/log/openwatch/worker.log"

. /etc/rc.subr

name="openwatch_worker"
rcvar="${name}_enable"

load_rc_config $name

: ${openwatch_worker_enable:="NO"}
: ${openwatch_worker_user:="openwatch"}
: ${openwatch_worker_logfile:="/var/log/openwatch/worker.log"}

pidfile="/var/run/${name}.pid"
command="/opt/openwatch/venv/bin/python3.12"
command_args="-m app.services.job_queue"

start_precmd="${name}_prestart"
stop_postcmd="${name}_poststop"

openwatch_worker_prestart()
{
    # Ensure log directory exists
    mkdir -p /var/log/openwatch
    chown "${openwatch_worker_user}" /var/log/openwatch

    # Set working directory and environment
    cd /opt/openwatch/backend || return 1
    export PYTHONPATH=/opt/openwatch/backend
    export PATH="/opt/openwatch/venv/bin:${PATH}"

    # Source environment file if it exists
    if [ -f /usr/local/etc/openwatch/secrets.env ]; then
        set -a
        . /usr/local/etc/openwatch/secrets.env
        set +a
    fi
}

openwatch_worker_poststop()
{
    rm -f "${pidfile}"
}

run_rc_command "$1"
