# packaging/deb/DEBIAN/prerm
#!/bin/bash
# OpenWatch Debian Pre-Removal Script
# Clean up before package removal

set -e

# Package name
PACKAGE="openwatch"

# Service management
stop_services() {
    echo "Stopping OpenWatch services..."
    
    # Stop services if they're running
    for service in openwatch openwatch-db; do
        if systemctl is-active --quiet $service.service 2>/dev/null; then
            systemctl stop $service.service || true
            echo "Stopped $service.service"
        fi
    done
    
    # Stop any remaining containers
    if command -v owadm >/dev/null 2>&1; then
        owadm stop --force 2>/dev/null || true
    fi
}

# Disable services
disable_services() {
    echo "Disabling OpenWatch services..."
    
    for service in openwatch openwatch-db; do
        if systemctl is-enabled --quiet $service.service 2>/dev/null; then
            systemctl disable $service.service || true
            echo "Disabled $service.service"
        fi
    done
}

# Remove AppArmor profile
remove_apparmor_profile() {
    if [ -f /etc/apparmor.d/openwatch-containers ]; then
        echo "Removing AppArmor profile..."
        
        # Unload the profile
        if command -v apparmor_parser >/dev/null 2>&1; then
            apparmor_parser -R /etc/apparmor.d/openwatch-containers 2>/dev/null || true
        fi
        
        # Remove the profile file
        rm -f /etc/apparmor.d/openwatch-containers
    fi
}

case "$1" in
    remove|upgrade|deconfigure)
        # Stop all services
        stop_services
        
        # Disable services on removal
        if [ "$1" = "remove" ]; then
            disable_services
        fi
        
        # Remove AppArmor profile
        remove_apparmor_profile
        
        # Note: We don't remove configuration or data directories here
        # They will be handled by postrm if this is a purge
        ;;
        
    failed-upgrade)
        ;;
        
    *)
        echo "prerm called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0