# packaging/deb/DEBIAN/prerm
#!/bin/bash
# OpenWatch Debian Pre-Removal Script
# Stop native services before package removal

set -e

case "$1" in
    remove|upgrade|deconfigure)
        echo "Stopping OpenWatch services..."

        # Stop services in dependency order
        for service in openwatch.target openwatch-beat.service \
                       openwatch-worker@1.service openwatch-api.service; do
            if systemctl is-active --quiet "$service" 2>/dev/null; then
                systemctl stop "$service" || true
                echo "  Stopped $service"
            fi
        done

        # Disable services on removal (not upgrade)
        if [ "$1" = "remove" ]; then
            echo "Disabling OpenWatch services..."
            for service in openwatch.target openwatch-beat.service \
                           openwatch-worker@1.service openwatch-api.service; do
                systemctl disable "$service" 2>/dev/null || true
            done
        fi

        # Remove AppArmor profile
        if [ -f /etc/apparmor.d/openwatch ]; then
            echo "Removing AppArmor profile..."
            if command -v apparmor_parser >/dev/null 2>&1; then
                apparmor_parser -R /etc/apparmor.d/openwatch 2>/dev/null || true
            fi
            rm -f /etc/apparmor.d/openwatch
        fi
        ;;

    failed-upgrade)
        ;;

    *)
        echo "prerm called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
